# Remote-Debug-App/supervisord.conf
[supervisord]
nodaemon=true

[program:api]
command=/usr/bin/dotnet /app/Remote-Debug-App.dll